# config/hermes/mm-agent.timer
[Unit]
Description=mm-agent run cadence (every 15 min)

[Timer]
OnBootSec=3min
OnUnitActiveSec=15min
RandomizedDelaySec=120
Persistent=true

[Install]
WantedBy=timers.target
